# DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable updatebkcmdb.service
systemctl start updatebkcmdb.service
